# slack/knockbot/knockbot.conf
start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5
umask 022

console none

env KNOCKBOT_TOKEN=xoxb-xxxxxxxxxxx-xxxxxxxxxxxxxxxxxxxxxxxx
env PB_USERNAME=xxxxx
env PB_PASSWORD=xxxxx

exec /usr/bin/nodejs /canadia/puzzlebitch/slack/knockbot/bin/bot.js
